# ziptronupdate.timer
[Unit]
Description=Запуск автодеплоя бота каждые 2 минуты

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s

[Install]
WantedBy=timers.target
